# rc.local
#!/bin/bash

bash /home/<YOURFILE>/kissdwaw.sh

exit 0
