# rc.local
#!/bin/bash

bash /home/<YOURFILE>/kiss-dw-aw.sh

exit 0
